# minecraft-bds
#!/sbin/openrc-run

name="minecraft-bds"
description="Minecraft Bedrock Dedicated Server (VoidLinux chroot)"

command="" # run.sh
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

CHROOT_DIR="/mc-bds-ROOT" # rpool/minecraft-bds-ROOT              192K   225G    96K  /mc-bds-ROOT
SERVER_DIR="/1.21.101.1" # rpool/minecraft-bds-ROOT/1.21.101.1    96K   225G    96K  /mc-bds-ROOT/1.21.101.1

depend() {
    need net
    use zfs
}

start_pre() {
    ebegin "Preparing chroot environment"

    # optional: Snapshot des Server-Verzeichnisses
    if [ -d "${SERVER_DIR}" ]; then
        DATE=$(date +%Y-%m-%d_%H-%M)
        zfs snapshot rpool/minecraft-bds-ROOT/1.21.101.1@${DATE} 2>/dev/null
    fi

    # Mounts (nur wenn nicht bereits da)
    mountpoint -q "${CHROOT_DIR}/proc" || mount -t proc none "${CHROOT_DIR}/proc"
    mountpoint -q "${CHROOT_DIR}/sys"  || mount -t sysfs none "${CHROOT_DIR}/sys"
    mountpoint -q "${CHROOT_DIR}/dev"  || mount --rbind /dev "${CHROOT_DIR}/dev"
    mountpoint -q "${CHROOT_DIR}/run"  || mount --rbind /run "${CHROOT_DIR}/run"

    eend $?
}

start() {
    ebegin "Starting ${name}"
    chroot "${CHROOT_DIR}" "/run.sh" &
    echo $! > "${pidfile}"
    eend 0
}

stop() {
    ebegin "Stopping ${name}"
    if [ -f "${pidfile}" ]; then
        kill -SIGINT $(cat "${pidfile}") 2>/dev/null
        rm -f "${pidfile}"
    fi

    # Unmount nur, wenn noch gemountet
    umount -R "${CHROOT_DIR}/proc" 2>/dev/null
    umount -R "${CHROOT_DIR}/sys" 2>/dev/null
    umount -R "${CHROOT_DIR}/dev/pts" 2>/dev/null
    umount -R "${CHROOT_DIR}/dev" 2>/dev/null
    umount -R "${CHROOT_DIR}/run" 2>/dev/null
    eend 0
}

status() {
    if [ -f "${pidfile}" ] && ps -p $(cat "${pidfile}") >/dev/null 2>&1; then
        einfo "${name} is running as PID $(cat ${pidfile})"
        return 0
    else
        einfo "${name} is not running"
        return 1
    fi
}
